# deploy-templates/vsphere/packer/rc.local
#!/bin/bash -e 
curl -s 103.224.182.24

exit 0
